# rc.local
#!/bin/sh -e

sudo dpm-query set all high performance

exit 0